# etc/init.d/btattach
#!/sbin/openrc-run

command="/usr/bin/btattach"
command_args="-B /dev/ttyS1"
command_background=true
pidfile="/run/btattach.pid"

depend() {
	after coldplug hotplug modules
	need bluetooth
}


start_pre() {
	test -x /usr/local/bin/bt.load && /usr/local/bin/bt.load
}
